# 1_server/autodeploy/webhook.ini
[program:webhook]
directory=/opt/POCKETCLOUD/SCRIPTS
autostart=true
autorestart=true
startretries=3
command=webhook -hooks /opt/hooks/hooks.json -logfile hooksout.log -hotreload